# 09-Apps/11-OracleXE/files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:oracle-xe]
command=/bin/bash /start-oracle-xe.sh

