# 08-Dev/11-java-base/files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:java]
command=/home/u01/workspace/start.sh